# systems/lctr-a72xx/fs/etc/systemd/system/cloud-init-local.service.d/mount-seed.conf
# Ensure our customized seed location is mounted prior to execution

[Unit]
RequiresMountsFor=/boot/firmware
